# docs/launchd/com.todotoday.frontend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.todotoday.frontend</string>
    <key>ProgramArguments</key>
    <array>
        <!-- UPDATE: absolute path to npx -->
        <string>/path/to/node/bin/npx</string>
        <string>vite</string>
        <string>--port</string>
        <string>5173</string>
    </array>
    <!-- UPDATE: absolute path to the frontend directory -->
    <key>WorkingDirectory</key>
    <string>/path/to/todo-today/frontend</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <!-- UPDATE: include your Node.js bin directory -->
        <string>/path/to/node/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/path/to/todo-today/logs/frontend.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/todo-today/logs/frontend.log</string>
</dict>
</plist>
